# etc/xdg/autostart/first-login-theme.desktop
[Desktop Entry]
Name=first-login-theme
Icon=
Exec=/usr/bin/first-login-theme
Terminal=false
Type=Application
